# launchd/com.appcaire.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.appcaire.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/node@22/bin/npx</string>
        <string>vite</string>
        <string>--port</string>
        <string>3010</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/be-agent-service/HomeCare/be-agents-service/apps/dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/be-agent-service/Library/Logs/appcaire-dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/be-agent-service/Library/Logs/appcaire-dashboard-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/opt/node@22/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
