# launchd/com.appcaire.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.appcaire.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/bjornevers_MacPro/HomeCare/be-agent-service/scripts/start-dashboard.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bjornevers_MacPro/HomeCare/be-agent-service</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/bjornevers_MacPro/Library/Logs/appcaire-dashboard.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bjornevers_MacPro/Library/Logs/appcaire-dashboard-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
