# launchd/com.appcaire.agent-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.appcaire.agent-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/bjornevers_MacPro/.nvm/versions/node/v22.15.1/bin/yarn</string>
        <string>workspace</string>
        <string>server</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/bjornevers_MacPro/HomeCare/be-agents-service</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/bjornevers_MacPro/.nvm/versions/node/v22.15.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PORT</key>
        <string>4010</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bjornevers_MacPro/Library/Logs/agent-server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/bjornevers_MacPro/Library/Logs/agent-server-error.log</string>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
